# puppet/modules/gov-test/files/hello_world.conf
start on runlevel [345]
exec /opt/hello_world/bin/gunicorn --bind 0.0.0.0:8000 --pythonpath /vagrant/hello_world hello_world:app